# DefenseLab/Services/SQLiSimple/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/bin/bash -c "/usr/sbin/mysqld"

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
